# devops/crontab
0 14 * * * /bin/bash /app/devops/bot.sh
